# systemd/mirror-world-marketing-knowledge-sync.timer
[Unit]
Description=Mirror World - Marketing knowledge sync every 2 hours

[Timer]
OnBootSec=15min
OnUnitActiveSec=2h
Persistent=true

[Install]
WantedBy=timers.target
